# Heedbook/ExpressHbApiTester/Scripts/apitester.service
[Unit]
Description=Express Api Tester for HbOnprem

[Service]
ExecStart=LAUNCHERPATH
Restart=always

[Install]
WantedBy=multi-user.target